# templates/etc/systemd/system/mailcatcher.service.erb
[Unit]
Description=Mailcatcher Service

[Service]
Restart=always
ExecStart=<%= @mailcatcher_path %>/mailcatcher<%= @options.join(' ') %> -f  >> /var/log/mailcatcher/mailcatcher.log 2>&1
ExecStop=/bin/true

[Install]
WantedBy=multi-user.target
