# templates/etc/systemd/system/mailcatcher.service.erb
[Unit]
Description=Mailcatcher Service
After=network.service

[Service]
Restart=always
ExecStart=<%= @mailcatcher_path %>/mailcatcher<%= @options.join(' ') %> --foreground

[Install]
WantedBy=multi-user.target

